# rc.d/rc.local
#!/bin/bash

/usr/bin/perl -T /opt/re-vend/boxStartup.pl 1>/var/log/revend/boxStartup 2>&1

exit 0;

